# systemd/eznap.service
[Unit]
Description="Trigger eznap"

[Service]
Type=oneshot
ExecStart=eznap trigger
