# systemd/eznap.service
[Unit]
Description="Trigger eznap"

[Service]
ExecStart=eznap trigger
