# systemd/eznap@.timer
[Unit]
Description=Trigger eznap on boot and then every %i minutes

[Timer]
OnBootSec=5min
OnActiveSec=%i

[Install]
WantedBy=timers.target
