# systemd/eznap-frequent.timer
[Unit]

[Timer]
OnUnitActiveSec=15m
Unit=eznap@frequent.service

[Install]
WantedBy=timers.target
