# systemd/eznap-hourly.timer
[Unit]

[Timer]
OnUnitActiveSec=1h
Unit=eznap@hourly.service

[Install]
WantedBy=timers.target
