# systemd/eznap-weekly.timer
[Unit]

[Timer]
OnUnitActiveSec=1w
Unit=eznap@weekly.service

[Install]
WantedBy=timers.target
