# systemd/eznap@.service
[Unit]
Description="Creates ZFS snapshots"

[Service]
ExecStart=eznap %i
