# usr/lib/systemd/system/eznap-daily.timer
[Unit]

[Timer]
OnUnitActiveSec=1d
Unit=eznap@daily

