# usr/lib/systemd/system/eznap-monthly.timer
[Unit]

[Timer]
OnUnitActiveSec=1M
Unit=eznap@monthly

